# etc/docker-scrapers.timer
[Unit]
Description=Run scrapers hourly

[Timer]
OnCalendar=hourly

[Install]
WantedBy=timers.target
