# etc/docker-notifier.timer
[Unit]
Description=Run notifier daily at 8AM

[Timer]
OnCalendar=08:00

[Install]
WantedBy=timers.target
